# supervisord.conf
[program:nginx]
command=nginx -g "daemon off;"
[supervisord]
nodaemon=true

[program:php-fpm]
command=php-fpm -F

